# src/skel/etc/supervisor.d/program-nginx.conf
; -*- ini -*-
[program:nginx]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
startretries = 100000000
autorestart = true
